# src/supervisord.conf
[program:elasticsearch]
command = /usr/local/bin/launch-elasticsearch.sh
